# environment/monitoring/monitor-smtp.timer
[Unit]
Description=Run monitor-smtp.service every 10 minutes

[Timer]
OnCalendar=*:0/10
